# root/etc/services.d/baidupcs/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid root /defaults/BaiduPCS-Go web --port $PORT